# scripts/systemd/artsfolio-backup.service
[Unit]
Description=ArtsFolio encrypted off-site backup
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=-/etc/artsfolio/artsfolio.env
EnvironmentFile=-/etc/artsfolio/backup.env
ExecStart=/var/www/artsfolio/scripts/backup/artsfolio_backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

# End of file.
