# scripts/systemd/artsfolio-email-worker@.service
[Unit]
Description=ArtsFolio email worker instance %i
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
WorkingDirectory=/var/www/artsfolio
Environment=ARTSFOLIO_ENV_FILE=/etc/artsfolio/artsfolio.env
Environment=ARTSFOLIO_WORKER_NAME=email-%i
Environment=ARTSFOLIO_EMAIL_STALE_MINUTES=30
ExecStart=/bin/bash /var/www/artsfolio/scripts/workers/email_loop.sh
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target

# End of file.
