# scripts/systemd/artsfolio-monitor.timer
[Unit]
Description=Run ArtsFolio operations health monitor every five minutes

[Timer]
OnCalendar=*-*-* *:0/5:00
AccuracySec=20s
Persistent=true
Unit=artsfolio-monitor.service

[Install]
WantedBy=timers.target
